# supervisor/systemd/aifishtank-claude-auth.service
[Unit]
Description=AI Fishtank Claude Auth Helper
Documentation=https://github.com/owner/ai-fishtank
After=network.target

[Service]
Type=simple
User=agent
Group=agent

ExecStartPre=+/bin/bash -c "mkdir -p /home/agent/.claude && chown -R agent:agent /home/agent/.claude && chmod 700 /home/agent/.claude"
ExecStart=/home/agent/ai-fishtank/supervisor/scripts/claude-auth-helper.sh

Restart=on-failure
RestartSec=5

WorkingDirectory=/home/agent/ai-fishtank

Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/agent/.local/bin:/home/agent/.npm-global/bin

# Resource constraints
MemoryMax=512M
TasksMax=128

# Security hardening
NoNewPrivileges=true
ProtectSystem=full

StandardOutput=journal
StandardError=journal
SyslogIdentifier=aifishtank-claude-auth

[Install]
WantedBy=multi-user.target
