# supervisor/systemd/aquarco-supervisor-python.service
[Unit]
Description=Aquarco Agent Supervisor (Python)
Documentation=https://github.com/owner/aquarco
After=aquarco-stack.service docker.service network.target
Wants=aquarco-stack.service

[Service]
Type=simple
User=agent
Group=agent

# Python entry point. Uses the installed package CLI.
ExecStart=/home/agent/.venv/bin/aquarco-supervisor run --config /home/agent/aquarco/supervisor/config/supervisor.yaml

# On SIGHUP the supervisor re-reads its config file without restarting.
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy: restart on unexpected failures, not on clean exit (code 0).
Restart=on-failure
RestartSec=10

# Give the process up to 60 seconds to drain in-flight tasks on SIGTERM.
TimeoutStopSec=60

WorkingDirectory=/home/agent/aquarco

# Standard PATH for finding gh, git, psql, claude, python, etc.
Environment=PATH=/home/agent/.venv/bin:/usr/local/bin:/usr/bin:/bin:/home/agent/.npm-global/bin
Environment=PYTHONUNBUFFERED=1

# Optional: load secrets from /etc/aquarco/secrets.env if it exists.
EnvironmentFile=-/etc/aquarco/secrets.env

# Resource constraints. Python uses slightly less memory than shell+subprocesses.
MemoryMax=2G
TasksMax=64

# Security hardening.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/agent /tmp

# Systemd-managed directories — created automatically on every service start,
# owned by User=/Group= (agent:agent). Survives /var/run tmpfs wipes on reboot.
RuntimeDirectory=aquarco
LogsDirectory=aquarco aquarco/agents
StateDirectory=aquarco

# Logging: structlog writes JSON to stderr; journalctl -u aquarco-supervisor-python follows.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aquarco-supervisor-python

[Install]
WantedBy=multi-user.target
